# etc directory files/systemd/system/motorcontrol_12h.timer
[Unit]
Description=Runs motor control script every 12 hours

[Timer]
OnBootSec=12h
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
