# etc directory files/systemd/system/motorcontrol_1h.timer
[Unit]
Description=Runs motor control script every 1 hour

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
